# php/crontab
*/5  *    *    *    *     phplist -pprocessqueue
45   3    *    *    *     phplist -pprocessbounces

